# daemons/com.local.paperless.ocr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.local.paperless.ocr</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/libexec/start-paperless-ocr.sh</string>
  </array>
  <!-- Runs as ROOT (no UserName), not TARGET_USER. macOS 15+/26 "Local Network Privacy"
       blocks non-Apple binaries (our venv python) from LAN access when the daemon runs as a
       user, causing "No route to host" reaching paperless-ngx — while /usr/bin/curl and ROOT
       daemons are exempt. Verified on macOS 26: Apple Vision (ocrmac/VNRecognizeTextRequest)
       works fine in a pure root daemon context, and root reaches the LAN. The wrapper still
       exports HOME/USER=TARGET_USER so caches/paths stay under the user's home. -->
  <key>WorkingDirectory</key><string>@TARGET_HOME@</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/var/log/macstudio/paperless-ocr.log</string>
  <key>StandardErrorPath</key><string>/var/log/macstudio/paperless-ocr.log</string>
</dict>
</plist>
