# daemons/com.local.immich.ml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.local.immich.ml</string>
  <key>UserName</key><string>@TARGET_USER@</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/libexec/start-immich-ml.sh</string>
  </array>
  <key>WorkingDirectory</key><string>@IMMICH_PROJECT_DIR@</string>
  <key>RunAtLoad</key><false/>
  <key>KeepAlive</key><false/>
  <key>ThrottleInterval</key><integer>5</integer>
  <key>StandardOutPath</key><string>/var/log/macstudio/immich-ml.log</string>
  <key>StandardErrorPath</key><string>/var/log/macstudio/immich-ml.log</string>
</dict>
</plist>
